# crontab
0 7 * * * cd /app && python crypto_monitor.py
